# examples/daemontools/6379.info/run
#!/bin/sh
php ../../../redis.info.php -i6379 -olocalhost -p6379 -w1