# examples/daemontools/6379.monitor/run
#!/bin/sh
php ../../../redis.monitor.php -i6379 -olocalhost -p6379 -t1 -c100 -w0